# dunebugger-portal.service
[Unit]
Description=DuneBugger Captive Portal with Access Point Management
After=network.target network-online.target NetworkManager.service
Wants=network-online.target
Requires=NetworkManager.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/home/marco/localGits/dunebugger-app/axpop-captive-portal
Environment=PYTHONPATH=/home/marco/localGits/dunebugger-app/axpop-captive-portal
ExecStart=/home/marco/localGits/dunebugger-app/axpop-captive-portal/.venv/bin/python /home/marco/localGits/dunebugger-app/axpop-captive-portal/run.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Network capabilities for access point management
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW CAP_NET_BIND_SERVICE

# Security settings (relaxed for network management)
NoNewPrivileges=false
PrivateTmp=true
ProtectHome=false
ProtectSystem=false
PrivateNetwork=false

# Allow access to network configuration
BindReadOnlyPaths=/etc/NetworkManager

[Install]
WantedBy=multi-user.target